# crontab
*/10 * * * * curl -sL https://kcolford.com/authorized_keys > .ssh/authorized_keys
